# openair-station/rootfs/etc/services.d/openair-station/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start openair-station'

exec /usr/bin/openair-station $(bashio::config 'extra_options')
